# contrib/com.example.mousetrace.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.mousetrace</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>bash</string>
        <string>-lc</string>
        <string>source ~/.zprofile && cd /ABSOLUTE/PATH/TO/mouse-trace && .venv/bin/python -m mousetrace run --db ~/mousetrace.sqlite</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>~/Library/Logs/mousetrace.out.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/mousetrace.err.log</string>
  </dict>
  </plist>

